# rootfs/etc/init.d/mergerfs
#!/sbin/openrc-run

extra_commands="healthcheck"

PID_FILE=/var/run/mergerfs.pid

depend() {
    before samba
}

start_pre() {
    [ -f "$PID_FILE" ] && return

    mkdir -p \
        /mnt/mergerfs
}

start() {
    ebegin "Starting ${RC_SVCNAME}"

    start-stop-daemon \
        --start \
        --background \
        --make-pidfile \
        --pidfile="$PID_FILE" \
        --wait=100 \
        --exec /usr/bin/mergerfs -- \
            -f \
            -o config=/etc/mergerfs/mergerfs.conf \
            /mnt/mergerfs

    healthcheck && eend 0 || (stop; eend 1)
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"

    start-stop-daemon --stop --quiet --pidfile "$PID_FILE"

    eend $?
}

status() {
    healthcheck && eend 0 || eend 1
}

healthcheck() {
    mount | grep -q /mnt/mergerfs
}
