# rootfs/etc/init.d/rclone
#!/sbin/openrc-run

extra_commands="healthcheck"

PID_FILE=/var/run/rclone.pid

depend() {
    before samba
}

start_pre() {
    [ -f "$PID_FILE" ] && return

    mkdir -p \
        /mnt/rclone \
        /var/rclone

    chown nobody \
        /mnt/rclone \
        /var/rclone

    echo "user_allow_other" >> /etc/fuse.conf
}

start() {
    ebegin "Starting ${RC_SVCNAME}"

    . /etc/rclone/.env

    start-stop-daemon \
        --start \
        --background \
        --make-pidfile \
        --pidfile="$PID_FILE" \
        --wait=100 \
        -u nobody \
        --exec /usr/bin/rclone -- \
        mount default: /mnt/rclone \
            --config=/etc/rclone/rclone.conf \
            --cache-dir=/tmp/rclone \
            --allow-other \
            --buffer-size=0 \
            --check-first \
            --checkers="$STREAM_LIMIT" \
            --cutoff-mode=soft \
            --multi-thread-chunk-size=1MiB \
            --multi-thread-cutoff=1MiB \
            --multi-thread-streams="$STREAM_LIMIT" \
            --no-checksum \
            --no-traverse \
            --no-update-dir-modtime \
            --no-update-modtime \
            --noappledouble=false \
            --server-side-across-configs \
            --transfers="$TRANSFER_LIMIT" \
            --use-mmap \
            --use-server-modtime \
            --vfs-cache-max-size="$VFS_CACHE_MAX_SIZE"MiB \
            --vfs-cache-mode=full \
            --vfs-disk-space-total-size=1PiB \
            --vfs-fast-fingerprint \
            --vfs-read-ahead=0 \
            --vfs-read-chunk-size=0 \
            --vfs-read-wait=0 \
            --vfs-write-back=1 \
            --vfs-write-wait=0

    healthcheck && eend 0 || (stop; eend 1)
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"

    start-stop-daemon --stop --quiet --pidfile "$PID_FILE"

    eend $?
}

status() {
    healthcheck && eend 0 || eend 1
}

healthcheck() {
    exit 0
    mount | grep -q "/mnt/rclone"
}
